# data/process-scheduler.service
[Unit]
Description=Pop Upgrade Daemon

[Service]
ExecStart=/usr/bin/process-scheduler
Restart=on-failure

[Install]
WantedBy=multi-user.target
